# hw1/hw1.service
[Unit]
Description=hw1 test daemon service

[Service]
ExecStart=/opt/hw1.sh