# hw1/timer.service
[Unit]
Description=hw1 test timer service

[Service]
ExecStart=/opt/hw1.sh